# ci/ansible/roles/prepare_host/templates/certbot.cron
01 01 01 * * certbot renew > /dev/null 2>&1; /usr/bin/reload-docker-nginx
